# server/lofi-stream-kick.service
[Unit]
Description=Lofi Stream to Kick (Arcade Theme)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/lofi-stream-kick/server
Environment=KICK_KEY=YOUR_STREAM_KEY_HERE
ExecStart=/opt/lofi-stream-kick/server/stream.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
